# etc/udev/crontab
@reboot /bin/sh -c '/home/nbleuzen/.dotfiles/etc/udev/wfh-setup.sh remove'
