# etc/profile.d/java_home.sh
# /etc/profile.d/java_home.sh
#
# Sets the location where ./bin/java can be found

export JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64
